# files/logship.cron
*/10 * * * * logshipper /usr/local/sbin/logship --config /etc/logship.conf --destination /srv/shared/logs
